# meta-jack/recipes-myDemo/mydemo1/files/hello.service
[Unit]
Description=GNU Hello World Startup script

[Service]
ExecStart=/home/root/hello-world
After=sample.service

[Install]
WantedBy=default.target
